# db/backup.crontab
*/5 * * * * /backup.sh >> /backup/backup.log
