# db/backup.crontab
*/2 * * * * /backup.sh >> /backup/backup.log 2>&1
